# systemd/ebook-bot.service
[Unit]
Description=PDF Ebook Translator Telegram Bot
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/ebook-translator
EnvironmentFile=/root/ebook-translator/.env
ExecStart=/root/ebook-translator/venv/bin/python /root/ebook-translator/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
