# systemd/ebook-translate@.service
[Unit]
Description=PDF Ebook Translator Translate Worker %i
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/ebook-translator
EnvironmentFile=/var/www/ebook-translator/.env
ExecStart=/var/www/ebook-translator/venv/bin/python /var/www/ebook-translator/translate_worker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

