# systemd/ebook-reconcile.service
[Unit]
Description=PDF Ebook Translator Reconciliation Worker
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/var/www/ebook-translator
EnvironmentFile=/var/www/ebook-translator/.env
ExecStart=/var/www/ebook-translator/venv/bin/python /var/www/ebook-translator/reconcile_worker.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
